# src/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:docker-gen]
command=docker-gen -only-exposed -watch -notify "supervisorctl restart haproxy" /usr/local/etc/haproxy/haproxy.tmpl /usr/local/etc/haproxy/haproxy.cfg

[program:haproxy]
command=haproxy -f /usr/local/etc/haproxy/haproxy.cfg
